# gateways/kong/files/kong.service
[Unit]
Description=Kong API Gateway
After=network.target
After=cassandra.service
Requires=cassandra.service

[Service]
ExecStart=/usr/local/bin/kong start
ExecStop=/usr/local/bin/kong stop
ExecReload=/usr/local/bin/kong reload
Type=forking
PIDFile=/usr/local/kong/kong.pid
User=kong
Group=kong

[Install]
WantedBy=default.target
