# gateways/apiman/files/apiman.service
[Unit]
Description=apiman
After=network.target

[Service]
WorkingDirectory=/usr/local/wildfly-9.0.2.Final
ExecStart=/usr/local/wildfly-9.0.2.Final/bin/standalone.sh -c standalone-apiman.xml
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=apiman
PrivateTmp=true
User=apiman
Group=apiman

[Install]
WantedBy=default.target
